# wificonnect.timer
[Unit]
Description=Run wificonnect.service every 2 minutes

[Timer]
OnBootSec=2min
OnCalendar=*:0/2

[Install]
WantedBy=multi-user.target